# systemd/system/cardano-db-sync@.service.d/override.conf
[Unit]
RequiresMountsFor=/blockchain

[Service]
BindPaths=/blockchain/cardano-db-sync/%i:/var/lib/cardano-db-sync/%i
