# systemd/system/cardano-node@.service.d/override.conf
[Unit]
RequiresMountsFor=/blockchain

[Service]
BindPaths=/blockchain/cardano/%i:/var/lib/cardano/%i
